# rootfs/etc/init.d/system-mgmt
#!/sbin/openrc-run
#
# OpenRC service script for System Management Web Application
#

name="system-mgmt"
description="System Management Web Application"

command="/usr/bin/python3"
command_args="/opt/system-mgmt/app.py"
command_background=true
command_user="root"

pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.log"

depend() {
    need net
    after firewall networking
}

start_pre() {
    # Ensure log directory exists
    checkpath --directory --mode 0755 /var/log

    # Ensure the app exists
    if [ ! -f /opt/system-mgmt/app.py ]; then
        eerror "System management app not found"
        return 1
    fi

    # Export conf.d variables so start-stop-daemon child inherits them
    export SYSTEM_MGMT_PORT="${SYSTEM_MGMT_PORT:-8000}"
    export SYSTEM_MGMT_HOST="${SYSTEM_MGMT_HOST:-0.0.0.0}"
}

start() {
    ebegin "Starting ${name}"
    start-stop-daemon --start \
        --make-pidfile \
        --pidfile "${pidfile}" \
        --background \
        --stdout "${output_log}" \
        --stderr "${error_log}" \
        --exec ${command} -- ${command_args}
    eend $?
}

stop() {
    ebegin "Stopping ${name}"
    start-stop-daemon --stop \
        --pidfile "${pidfile}" \
        --retry 10
    eend $?
}

status() {
    if [ -f "${pidfile}" ]; then
        if kill -0 $(cat "${pidfile}") 2>/dev/null; then
            einfo "${name} is running"
            return 0
        fi
    fi
    einfo "${name} is not running"
    return 1
}
